# data/example.calendar-crawler.timer
[Unit]
Description=Run Calendar Crawler daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
